# contrib/init/bellsd.conf
description "LinkCoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env linkcoind_BIN="/usr/bin/linkcoind"
env linkcoind_USER="linkcoin"
env linkcoind_GROUP="linkcoin"
env linkcoind_PIDDIR="/var/run/linkcoind"
# upstart can't handle variables constructed with other variables
env linkcoind_PIDFILE="/var/run/linkcoind/linkcoind.pid"
env linkcoind_CONFIGFILE="/etc/linkcoin/linkcoin.conf"
env linkcoind_DATADIR="/var/lib/linkcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # linkcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$linkcoind_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run linkcoind."
        echo "The setting must appear in $linkcoind_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"LinkCoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$linkcoind_PIDDIR"
    chmod 0755 "$linkcoind_PIDDIR"
    chown $linkcoind_USER:$linkcoind_GROUP "$linkcoind_PIDDIR"
    chown $linkcoind_USER:$linkcoind_GROUP "$linkcoind_CONFIGFILE"
    chmod 0660 "$linkcoind_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$linkcoind_PIDFILE" \
    --chuid $linkcoind_USER:$linkcoind_GROUP \
    --exec "$linkcoind_BIN" \
    -- \
    -pid="$linkcoind_PIDFILE" \
    -conf="$linkcoind_CONFIGFILE" \
    -datadir="$linkcoind_DATADIR" \
    -disablewallet \
    -daemon

